# ansible/roles/brain/templates/brain.service.j2
[Unit]
Description=Brain
After=network.target
[Service]
ExecStart=/home/brain/brain
WorkingDirectory=/home/brain
Restart=always
User=brain
Group=brain
Environment=PATH=/usr/bin:/usr/local/bin